# debian-pkg/typesense-server/DEBIAN/postinst
#!/bin/sh

systemctl enable typesense-server
systemctl daemon-reload
systemctl start typesense-server